# src/vsftpd/build/modules/build/files/etc/supervisor/conf.d/vsftpd.conf
[program:vsftpd]
command = /usr/sbin/vsftpd
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
